# templates/camviewport-reload.service.j2
[Unit]
Description=Reload camviewport

[Service]
ExecStart={{ camviewport_home }}/reload.sh
User={{ camviewport_user }}
Group={{ camviewport_group }}
Type=oneshot
